# systemd/system/snap.docker.dockerd.service
[Unit]
# Auto-generated, DO NOT EDIT
Description=Service for snap application docker.dockerd
Requires=snap-docker-796.mount
Wants=network.target
After=snap-docker-796.mount network.target snapd.apparmor.service
X-Snappy=yes

[Service]
EnvironmentFile=-/etc/environment
ExecStart=/usr/bin/snap run docker.dockerd
SyslogIdentifier=docker.dockerd
Restart=on-failure
WorkingDirectory=/var/snap/docker/796
TimeoutStopSec=30
Type=simple
Delegate=true

[Install]
WantedBy=multi-user.target
