# systemd/system/snap-google\x2dcloud\x2dsdk-196.mount
[Unit]
Description=Mount unit for google-cloud-sdk, revision 196
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/google-cloud-sdk_196.snap
Where=/snap/google-cloud-sdk/196
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
